# deploy/systemd/smartlife-primary-mqtt.service
[Unit]
Description=SmartLife Primary MQTT Broker on 19283
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/smartlife-primary
ExecStart=/usr/sbin/mosquitto -c /home/ubuntu/smartlife-primary/deploy/mosquitto.conf
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
